# deploy/legacy-systemd/betterblackboard.service
[Unit]
Description=BetterBlackboard web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bb
WorkingDirectory=/opt/betterblackboard
EnvironmentFile=/opt/betterblackboard/.env
ExecStart=/opt/betterblackboard/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
